# package/softwinner/wifimanager/udhcpc.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=22
STOP=99

USE_PROCD=1
PROG=/sbin/udhcpc

start_service() {
    procd_open_instance 
    procd_set_param command $PROG -i wlan0 
    procd_close_instance 
}

stop_service() {
    procd_open_instance
    killall -q -KILL udhcpc 
    procd_close_instance
}

restart(){
    stop "$@"
    start "$@"
}